# interactivecv.service
[Unit]
Description=Interactive CV — Flask + SocketIO
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/interactiveCV
EnvironmentFile=/home/ubuntu/interactiveCV/.env

# gunicorn with eventlet worker is required for SocketIO WebSocket support.
# -w 1 is mandatory with eventlet — multiple workers break SocketIO's
# in-memory state (admin_state dict, room routing).
ExecStart=/home/ubuntu/interactiveCV/venv/bin/gunicorn \
    --worker-class eventlet \
    -w 1 \
    --bind 0.0.0.0:5000 \
    --timeout 120 \
    --log-level info \
    main:app

Restart=always
RestartSec=5

# Re-register the Telegram webhook on every start using the current
# public IP. This handles the case where the EC2 instance was restarted
# and got a new IP (no static/Elastic IP assigned).
ExecStartPost=/home/ubuntu/interactiveCV/register_webhook.sh

[Install]
WantedBy=multi-user.target
